# docker/cron.d/mail
* * * * * cd /helpy && bundle exec rake helpy:mailman
