# docker/cron.d/mail
* * * * * bundle exec rake helpy:mailman --silent
