# docker/cron.d/mail
@reboot /bin/bash -l -c 'RAILS_ENV=production bundle exec rake helpy:mailman mail_interval=60'
